# reranker/reranker.service
[Unit]
Description=BGE Reranker Service
After=network.target

[Service]
Type=simple
User=povniouk
WorkingDirectory=/home/povniouk/reranker
ExecStart=/home/povniouk/reranker/venv/bin/python app.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
